# etc/services.d/sabnzbd/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc \
  /usr/bin/python3 -OO /sabnzbd/SABnzbd.py -f ${SABNZBD_HOME} \
  -s ${SABNZBD_BIND_ADDRESS}:${SABNZBD_PORT} ${SABNZBD_OPTS}

